# elm327-visdatafeeder/systemd/elm327-visdatafeeder.service
[Unit]
Description = OBDII data feeder thLM 327 for w3c-visserver
After=w3c-visserver.service
 
[Service]
WorkingDirectory=/usr/bin/elm327-visdatafeeder
ExecStart=/usr/bin/elm327-visdatafeeder/start.sh

[Install]
WantedBy=multi-user.target
